# src/integration/vagrant/etcd.conf
start on vagrant-mounted
stop on runlevel [016]

respawn
respawn limit 10 5

script
/vagrant/src/main/bash/etcd_run.sh -b 0.0.0.0 /opt/etcd-v2.2.5-linux-amd64 \
  >>/var/log/etcd.log 2>&1
end script
